# vlc/root/etc/services.d/vlc/run
#!/usr/bin/with-contenv bashio

CONFIG_PATH=/data/options.json
TELNET_PASS="$(jq --raw-output '.password' $CONFIG_PATH)"

cvlc --version

bashio::log.info 'Start vlc'

cvlc -I telnet --telnet-password $TELNET_PASS --telnet-port 4212 --aout=pulse 2>/dev/null
